# config/goxel-daemon.service
[Unit]
Description=Goxel 3D Voxel Editor Daemon
Documentation=https://goxel.xyz
After=network.target
Wants=network.target

[Service]
Type=forking
User=goxel
Group=goxel
WorkingDirectory=/var/lib/goxel
ExecStartPre=/usr/bin/test -f /usr/local/bin/goxel-daemon
ExecStart=/usr/local/bin/goxel-daemon --daemonize --pid-file /var/run/goxel-daemon.pid --socket /var/run/goxel-daemon.sock --log-file /var/log/goxel-daemon.log --workers 4 --queue-size 1024 --max-connections 256
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
ExecStopPost=/bin/rm -f /var/run/goxel-daemon.sock
PIDFile=/var/run/goxel-daemon.pid
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=goxel-daemon

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/goxel /var/log /var/run
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=false
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Health monitoring
TimeoutStartSec=30
TimeoutStopSec=30
WatchdogSec=60

[Install]
WantedBy=multi-user.target