# config/systemd/goxel-daemon.service
[Unit]
Description=Goxel 3D Voxel Editor Daemon
Documentation=man:goxel-daemon(8)
After=network.target
Wants=network.target

[Service]
Type=notify
User=goxel
Group=goxel
Environment=GOXEL_HEADLESS=1
Environment=OSMESA_RENDERING=1

# Main service configuration
ExecStart=/usr/local/bin/goxel-daemon \
    --foreground \
    --socket /var/run/goxel/daemon.sock \
    --pid-file /var/run/goxel/daemon.pid \
    --log-file /var/log/goxel/daemon.log \
    --workers 4 \
    --max-connections 256

# Pre-start setup
ExecStartPre=/bin/mkdir -p /var/run/goxel
ExecStartPre=/bin/mkdir -p /var/log/goxel
ExecStartPre=/bin/chown goxel:goxel /var/run/goxel
ExecStartPre=/bin/chown goxel:goxel /var/log/goxel

# Post-stop cleanup
ExecStopPost=/bin/rm -f /var/run/goxel/daemon.sock
ExecStopPost=/bin/rm -f /var/run/goxel/daemon.pid

# Resource management
LimitNOFILE=65536
LimitNPROC=32768

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/run/goxel /var/log/goxel /tmp

# Restart policy
Restart=always
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=30

# Service management
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
Alias=goxel.service